# tinyirc.service
[Unit]
Description=TinyIRC Service

[Service]
Type=simple
ExecStart=/usr/local/bin/tinyircd

[Install]
WantedBy=multi-user.target